# ops/launchd/tradenest-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.tradenest.dashboard</string>

  <!--
    Runtime working directory. Wrapper script will cd into dashboard/
    after loading .env from the runtime root.
  -->
  <key>WorkingDirectory</key>
  <string>/Users/bharatmacmini/tradenest-runtime/dashboard</string>

  <!--
    Wrapper script handles cd, .env loading, NEXT_PUBLIC_TRADENEST_API_BASE
    fallback, timestamped start log, and exec of npm so launchd tracks the
    real Node PID. /bin/zsh -lc gives us the user's login PATH (nvm, etc.).
  -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>exec /Users/bharatmacmini/tradenest-runtime/ops/scripts/run-dashboard.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/bharatmacmini/tradenest-runtime/logs/dashboard.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bharatmacmini/tradenest-runtime/logs/dashboard.err.log</string>
</dict>
</plist>
