# ops/launchd/tradenest-backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.tradenest.backend</string>

  <key>WorkingDirectory</key>
  <string>/Users/bharatmacmini/Documents/New project</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>source .env &amp;&amp; .venv/bin/python -m uvicorn tradenest.main:app --host 127.0.0.1 --port 8000</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/bharatmacmini/Documents/New project/logs/backend.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bharatmacmini/Documents/New project/logs/backend.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
