# ops/launchd/tradenest-backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.tradenest.backend</string>

  <!--
    Runtime working directory. The repo (source of truth) lives elsewhere;
    this directory is the deployed runtime copy that macOS launchd is
    permitted to read. See docs/runbook.md for repo vs runtime layout.
  -->
  <key>WorkingDirectory</key>
  <string>/Users/bharatmacmini/tradenest-runtime</string>

  <!--
    Wrapper script handles cd, .env loading, timestamped start log, and
    exec of the real backend process so launchd tracks the right PID.
    /bin/zsh -lc gives us the user's login PATH (homebrew, pyenv, etc.).
  -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>exec /Users/bharatmacmini/tradenest-runtime/ops/scripts/run-backend.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/bharatmacmini/tradenest-runtime/logs/backend.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bharatmacmini/tradenest-runtime/logs/backend.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
